# cron.d/mediawiki2fedora.cron
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin

0 0 * * * lludss /opt/mediawiki2fedora/bin/mediawiki2fedora 1>> /var/log/mediawiki2fedora/mediawiki2fedora.log 2>> /var/log/mediawiki2fedora/mediawiki2fedora.error
